# images/php7-fpm-trusty64/configs/supervisord/supervisord.conf
[supervisord]
nodaemon = true
[program:php-fpm7.0]
command = /usr/sbin/php-fpm7.0
user = root
autostart = true